# crontab
# TODO: calculate sunrise/sunset and run at jobs for these
0 8 * * * mocp -p
0 21 * * * mocp -s
